# systemd/system/clumsy.service
[Unit]
Description=My little clumsy test daemon

[Service]
Type=simple
ExecStart=/usr/bin/false

[Install]
WantedBy=multi-user.target
