# systemd/system/instanced.socket
[Unit]
Description=instanced socket

[Socket]
ListenStream=50100
Accept=True

[Install]
WantedBy=sockets.target