# systemd/system/my.service
[Unit]
Description=My oneshot unit
Requires=required.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/bin/usleep 1


