# systemd/system/my_little_daemon.service
[Unit]
Description=My little test daemon

[Service]
Type=simple
ExecStart=/usr/bin/my_little_daemon

[Install]
WantedBy=multi-user.target
