# systemd/system/oneshot.service
[Unit]
Description=Simple oneshot unit

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/oneshot.sh hello

