# systemd/system/really-required.service
[Unit]
Description=My oneshot unit
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/bin/usleep 1

