# systemd/system/screamer.service
[Unit]
Description=Screamer's oneshot unit

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/screamer
